# sysctl-logger.service.in
[Unit]
Description=sysctl monitoring with BPF
DefaultDependencies=no
Before=systemd-sysctl.service

[Service]
ExecStart=${SBINDIR}/sysctl-logger
StandardError=journal
StandardOutput=journal
StandardInput=null
LimitMEMLOCK=infinity # Might not be necessary

[Install]
WantedBy=basic.target
