# sysctl-logger.service.in
[Unit]
Description=sysctl monitoring with BPF
DefaultDependencies=no
Wants=systemd-sysctl.service
Before=systemd-sysctl.service

[Service]
Type=notify
ExecStart=${SBINDIR}/sysctl-logger
StandardError=journal
StandardOutput=journal
StandardInput=null
LimitMEMLOCK=infinity # Might not be necessary

[Install]
WantedBy=sysinit.target
